# files/r7800-router/etc/init.d/irq_cpu
#!/bin/sh /etc/rc.common

START=15

boot() {
	echo 2 > /proc/irq/39/smp_affinity
	echo 2 > /proc/irq/40/smp_affinity
}
