# files/rpi-router/etc/init.d/eth_cpu
#!/bin/sh /etc/rc.common

START=15

boot() {
	echo c > /proc/irq/35/smp_affinity
	echo 2 > /proc/irq/36/smp_affinity
}
