# files/rpi5-router/etc/init.d/eth_cpu
#!/bin/sh /etc/rc.common

START=15

boot() {
	echo 8 > /sys/class/net/eth1/queues/rx-0/rps_cpus
	echo 4 > /sys/class/net/eth0/queues/rx-0/rps_cpus
}
