# files/rpi5-router/etc/init.d/eth_cpu
#!/bin/sh /etc/rc.common

START=15

boot() {
	echo 8 > /proc/irq/103/smp_affinity
	echo 4 > /proc/irq/133/smp_affinity
	echo 2 > /proc/irq/128/smp_affinity
}
